# lib/scripts/ddv
#!/sbin/runscript

extra_started_commands="reload"

DDV=%DDV_PATH%
USER=%USER%

export PATH=%NODE_PATH%:$PATH
export DDV_HOME_PATH="%DDV_HOME_PATH%"
export DDV_CONFIG_PATH="%DDV_CONFIG_PATH%"
export DDV_OUT_LOG_FILE_PATH="%DDV_OUT_LOG_FILE_PATH%"
export DDV_ERR_LOG_FILE_PATH="%DDV_ERR_LOG_FILE_PATH%"
export HOME_PATH="%HOME_PATH%"

super() {
    su - $USER -c "PATH=$PATH; DDV_HOME_PATH=$DDV_HOME_PATH DDV_CONFIG_PATH=$DDV_CONFIG_PATH DDV_OUT_LOG_FILE_PATH=$DDV_OUT_LOG_FILE_PATH DDV_ERR_LOG_FILE_PATH=$DDV_ERR_LOG_FILE_PATH HOME_PATH=$HOME_PATH $*"
}

depend() {
	use net localmount
	after bootmisc
}

start() {
	ebegin "Starting DDV"

	super $DDV start

	eend $?
}

stop() {
	ebegin "Stopping DDV"

	#super $DDV dump
	super $DDV stop
	super $DDV kill

	eend
}

reload() {
	ebegin "Reloading DDV"

	super $DDV reload all

	eend $?
}

status() {
	ebegin "Status for DDV"

	super $DDV list

	eend $?
}
